# infra/systemd/preenrich-next-stage-sweeper.service
[Unit]
Description=Iteration 4 Preenrich Next-Stage Sweeper
After=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/root/scout-cron
EnvironmentFile=/root/scout-cron/.env
ExecStart=/root/scout-cron/.venv/bin/python -m src.preenrich.sweepers next-stage --limit 100
