# infra/systemd/preenrich-cv-ready-sweeper.service
[Unit]
Description=Iteration 4 Preenrich CV Ready Finalizer
After=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/root/scout-cron
EnvironmentFile=/root/scout-cron/.env
ExecStart=/root/scout-cron/.venv/bin/python -m src.preenrich.sweepers cv-ready --limit 100
