# infra/systemd/preenrich-stage-sweeper.timer
[Unit]
Description=Run iteration 4 preenrich stage lease sweeper every 60 seconds

[Timer]
OnBootSec=60s
OnUnitActiveSec=60s
AccuracySec=5s
Unit=preenrich-stage-sweeper.service

[Install]
WantedBy=timers.target
